# deploy/systemd/lupa-na-prawo.service
[Unit]
Description=Lupa na prawo (Web UI + API + MCP)
After=network.target postgresql.service ollama.service

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/lupa-na-prawo
EnvironmentFile=/home/ec2-user/lupa-na-prawo/.env
ExecStart=/home/ec2-user/lupa-na-prawo/.venv/bin/python -m uvicorn polish_law_helper.server:app --host 127.0.0.1 --port 8766
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
